# rc_scripts/iredapd.openbsd
#!/bin/sh

daemon='python /opt/iredapd/iredapd.py'

. /etc/rc.d/rc.subr

rc_reload=NO
rc_cmd $1
